# standalone/chef-repo/cookbooks/concourse/templates/etc/systemd/system/concourse-worker.service.erb
[Unit]
Description=concourse-worker
After=network.target

[Service]
User=root
ExecStart=/usr/local/concourse/bin/concourse worker --work-dir /opt/concourse/worker --tsa-host <%= node['concourse_worker']['tsa_host'] %>  --tsa-public-key /usr/local/concourse/keys/tsa_host_key.pub --tsa-worker-private-key /usr/local/concourse/keys/worker_key

[Install]
WantedBy=multi-user.target
