# standalone/chef/cookbooks/concourse/templates/etc/systemd/system/concourse-web.service.erb
[Unit]
Description=concourse-web
After=network.target

[Service]
User=concourse
ExecStart=/usr/local/concourse/bin/concourse web --basic-auth-username <%= node['concourse_web']['basic_auth_username'] %> --basic-auth-password <%= node['concourse_web']['basic_auth_password'] %> --session-signing-key /usr/local/concourse/keys/session_signing_key --tsa-host-key /usr/local/concourse/keys/tsa_host_key --tsa-authorized-keys /usr/local/concourse/keys/authorized_worker_keys --postgres-data-source postgres://concourse@localhost:5432/concourse?sslmode=disable --external-url <%= node['concourse_web']['external_url'] %>

[Install]
WantedBy=multi-user.target
